# templates/etc/systemd/system/nftables.service.d/override.conf.j2
# {{ ansible_managed }}

[Service]
{% if not nft__service_protect %}
ProtectSystem=no
ProtectHome=no
{% endif %}
