# templates/etc/systemd/system/fail2ban.service.d/override.conf.j2
{{ ansible_managed | comment }}

[Unit]
After=network.target iptables.service firewalld.service ip6tables.service ipset.service nftables.service
PartOf=firewalld.service nftables.service

[Install]
WantedBy=multi-user.target nftables.service
